# pimdd.service.in
[Unit]
Description=PIM-DM multicast routing daemon
Documentation=man:pimdd
Documentation=file:@DOCDIR@/README.md
Documentation=https://github.com/troglobit/pimdd
# ConditionPathExists=@SYSCONFDIR@/pimdd.conf
After=network-online.target
Requires=network-online.target

[Service]
Type=simple
ExecStart=@SBINDIR@/pimdd -ns

[Install]
WantedBy=multi-user.target
