# pimdd.service.in
[Unit]
Description=PIM-DM multicast routing daemon
Documentation=man:pimdd
Documentation=man:pimctl
Documentation=man:pimdd.conf
# ConditionPathExists=@SYSCONFDIR@/pimdd.conf
After=network-online.target
Requires=network-online.target

[Service]
Type=simple
EnvironmentFile=-@SYSCONFDIR@/default/pimd
ExecStart=@SBINDIR@/pimdd -ns $PIMDD_OPTS

[Install]
WantedBy=multi-user.target
